# bitlbee/root/etc/services.d/bitlbee/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/bin/bitlbee \
    -P /tmp/bitlbee.pid -F -n -c /config/bitlbee.conf -d /config/user-conf/
