# fatso/root/etc/crontabs/root
# run hourly backup rotate
*/1 * * * * /update-healthcheck.sh
